# systemd-examples/feeds-to-pocket.service
[Unit]
Description=Feeds to Pocket

[Service]
Type=oneshot
ExecStart=feeds-to-pocket ~/feeds-to-pocket.yaml
