# systemd-examples/feeds-to-pocket.timer
[Unit]
Description=Timer for Feeds to Pocket

[Timer]
OnStartupSec=1m
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
